# crontab/root
*/5 * * * * /system/sdcard/scripts/checkwifi.sh > /dev/null
